# branches/1.2/debian/jobmonarch-jobarchived.prerm
#!/bin/sh

PATH=/bin:/sbin:/usr/bin:/usr/sbin

set -e

case "$1" in
    deconfigure)
        /etc/init.d/jobmonarch-jobarchived stop || true
        ;;
    remove)
        # Only on removal (not during upgrade)
        /etc/init.d/jobmonarch-jobarchived stop || true
        update-rc.d -f jobmonarch-jobarchived remove || true
        ;;
    *)
        ;;
esac

exit 0
