# branches/1.0/debian/jobmonarch-jobarchived.prerm
#!/bin/sh

PATH=/bin:/sbin:/usr/bin:/usr/sbin

set -e

case "$1" in
    deconfigure)
        /etc/init.d/jobarchived stop || true
        ;;
    remove)
        # Only on removal (not during upgrade)
        /etc/init.d/jobarchived stop || true
        update-rc.d -f jobarchived remove || true
        ;;
    *)
        ;;
esac

exit 0
