# pkg/deb/jobmond/DEBIAN/postinst
#!/bin/sh

PATH=/bin:/sbin:/usr/bin:/usr/sbin

if [ -x /etc/init.d/jobmond ] && [ -x /usr/sbin/jobmond ]
	then

		update-rc.d -f jobmond defaults

fi

/etc/init.d/jobmond restart
