# pkg/deb/jobmond/DEBIAN/postrm
#!/bin/sh

PATH=/bin:/sbin:/usr/bin:/usr/sbin

/etc/init.d/jobmond stop
update-rc.d -f jobmond remove
