# debian/jobmonarch-jobmond.prerm
#!/bin/sh

PATH=/bin:/sbin:/usr/bin:/usr/sbin

set -e

case "$1" in
    remove|deconfigure)
        # Only on removal (not during upgrade)
        /etc/init.d/jobmond stop || true
        update-rc.d -f jobmond remove || true
        ;;
    *)
        ;;
esac

exit 0
